# com.jeffzeng.whale-scanner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jeffzeng.whale-scanner</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/jeffzeng/Projects/polymarket-whale-scanner/dist/index.js</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>WorkingDirectory</key>
    <string>/Users/jeffzeng/Projects/polymarket-whale-scanner</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jeffzeng/Projects/polymarket-whale-scanner/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jeffzeng/Projects/polymarket-whale-scanner/logs/stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
